# services/toggle_theme.timer
[Unit]
Description=Timer para ejecutar toggle_theme

[Timer]
OnCalendar=06:01
OnCalendar=17:01
Persistent=true

[Install]
WantedBy=timers.target

